# docker/data/service/peaks-unpack/run
#!/bin/sh
peaks --stdout unpack --thread 1
