# docker/data/service/peaks-unpack/run
#!/bin/sh
peaks --stdout unpack --threads 1
